# ansible/roles/openclaw/templates/openclaw-gateway.service.j2
[Unit]
Description=OpenClaw Gateway
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw
EnvironmentFile=/etc/openclaw/gateway.env
ExecStart=/usr/bin/openclaw gateway --port 18789 --verbose
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
