# ansible/roles/openclaw/templates/openclaw-gateway.service.j2
[Unit]
Description=OpenClaw Gateway
After=network-online.target docker.service
Wants=network-online.target docker.service

[Service]
Type=simple
User=openclaw
Group=openclaw
EnvironmentFile=/etc/openclaw/gateway.env
Environment=PATH=/home/openclaw/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/home/openclaw/.local/bin/openclaw gateway --port 18789 --verbose
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
PrivateTmp=true
ReadWritePaths=/home/openclaw/.openclaw /home/openclaw/.local

[Install]
WantedBy=multi-user.target
