# ansible/roles/openclaw/templates/openclaw-alert@.service.j2
# {{ ansible_managed }}
[Unit]
Description=Send Telegram alert for failed unit %i

[Service]
Type=oneshot
ExecStart=/usr/local/bin/openclaw-telegram-alert.sh %i
