# ansible/roles/openclaw/templates/openclaw-gateway.service.j2
[Unit]
Description=OpenClaw Gateway
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw
{% if llm_api_key != "" %}
{% if llm_provider == "anthropic" %}
Environment="ANTHROPIC_API_KEY={{ llm_api_key }}"
{% elif llm_provider == "openai" %}
Environment="OPENAI_API_KEY={{ llm_api_key }}"
{% elif llm_provider == "openrouter" %}
Environment="OPENROUTER_API_KEY={{ llm_api_key }}"
{% elif llm_provider == "gemini" %}
Environment="GEMINI_API_KEY={{ llm_api_key }}"
{% elif llm_provider == "do_ai" %}
Environment="DO_AI_API_KEY={{ llm_api_key }}"
{% elif llm_provider == "gradient" %}
Environment="GRADIENT_API_KEY={{ llm_api_key }}"
{% else %}
Environment="{{ llm_env_var | default('LLM_API_KEY') }}={{ llm_api_key }}"
{% endif %}
{% endif %}
ExecStart=/usr/bin/openclaw gateway --port 18789 --verbose
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
