# ansible/roles/backup/templates/openclaw-backup.service.j2
[Unit]
Description=OpenClaw Restic Backup
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/bin/docker compose -f /opt/openclaw/docker-compose.yml run --rm restic restic backup /data --exclude='*.log' --exclude='sandboxes/*'
ExecStartPost=/usr/bin/docker compose -f /opt/openclaw/docker-compose.yml run --rm restic restic forget --keep-hourly 24 --keep-daily 7 --keep-weekly 4 --prune
