# ansible/roles/backup/templates/openclaw-backup-check.service.j2
[Unit]
Description=Verify OpenClaw Restic Backup Integrity
After=network-online.target
Wants=network-online.target
{% if telegram_bot_token != "" and telegram_alert_chat_id != "" %}
OnFailure=openclaw-alert@%n.service
{% endif %}

[Service]
Type=oneshot
User=openclaw
EnvironmentFile=/home/openclaw/.restic-env
ExecStart=/usr/bin/restic check
